# server/aas-password-server.service
[Unit]
Description=AAS Password Authentication Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aasauth
Group=aasauth
WorkingDirectory=/opt/aas-password-server
EnvironmentFile=/etc/aas-password-server.env
ExecStart=/opt/aas-password-server/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
